# overlay/etc/periodic/15min/fr-cron
#!/bin/sh

set -euo pipefail


#======================================================================================================================
# Call freshrss-refresh executable
#======================================================================================================================

fr-cron
