# overlay/etc/periodic/15min/fr-cron
#!/bin/sh

set -euo pipefail


#======================================================================================================================
# Call fr-actualize executable
#======================================================================================================================

fr-actualize
